# btdht.service
[Unit]
Description=BT DHT
After=network.target

[Service]
ExecStart=%h/sources/btdht/dht-server -b 11000
Restart=always
